# deploy/roles/monitor_hyper/templates/packet-drop-exporter.service.j2
[Unit]
Description=Export packet drop indicators from dmesg to Prometheus metrics
After=network.target

[Service]
Type=oneshot
User=root
ExecStart=/opt/cloudland/scripts/monitor/export_packet_drops.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
